# app/vps-scripts/wantokjobs.service
[Unit]
Description=WantokJobs Web Application
Documentation=https://wantokjobs.com
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/wantokjobs/app/server
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=10

# Environment
Environment=NODE_ENV=production
Environment=PORT=3001

# Load additional environment from .env file
EnvironmentFile=-/opt/wantokjobs/app/server/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wantokjobs

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
TimeoutStartSec=90s
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target
